# releases/green/ssop/etc/systemd/system/ssop_gunicorn.socket
[Unit]
Description=Single Sign One Portal gunicorn socket

[Socket]
ListenStream=/run/ssop_gunicorn.sock

[Install]
WantedBy=sockets.target
